# conf/lbbs.service.default
[Unit]
Description=LBBS Server
After=network-online.target

[Service]
User=bbs
Group=bbs

Type=notify
ExecStart=/usr/local/lbbs/bin/bbsd -f
KillSignal=SIGTERM
ExecReload=/bin/kill -HUP $MAINPID

[Install]
WantedBy=multi-user.target
